# files/backup.service
[Unit]
Description=Database backup service
[Service]
Type=oneshot
ExecStart=/root/db_backup.sh